# crontab
*/1 * * * * cd /app/src && /usr/local/bin/python3 run_tg_grabber.py > /app/logs/cron_log.log 2>/app/logs/cron_err.log
